# conf/supervisor/fingerweb.conf
[program:fingerweb]
command=bash -c "./manage.py migrate && gunicorn fingerweb.wsgi:application --bind 0.0.0.0:8000 --workers 3"
autorestart=true
startsecs=3
user=finger
directory=/app
stdout_events_enabled=true
stderr_events_enabled=true
