# conf/supervisor/nginx.conf
[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autorestart=true
startsecs=3
stdout_events_enabled=true
stderr_events_enabled=true
